# config/sway/scripts/enable-iio-buffer.service
[Unit]
Description=Enable IIO accelerometer buffer
After=systemd-modules-load.service

[Service]
Type=oneshot
# Retry enabling the buffer for up to 60 seconds — the device may not be ready at boot.
ExecStart=/bin/sh -c 'for i in $(seq 1 20); do \
    if [ -f /sys/bus/iio/devices/iio:device0/buffer/enable ]; then \
        echo 1 > /sys/bus/iio/devices/iio:device0/buffer/enable 2>/dev/null && exit 0; \
    fi; \
    sleep 3; \
done; \
echo "Failed to enable IIO buffer after 60s" >&2; exit 1'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
