# config/sway/scripts/enable-iio-buffer.service
[Unit]
Description=Keep IIO accelerometer buffer enabled
After=systemd-modules-load.service

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/bin/sh -c 'while true; do \
    if [ -f /sys/bus/iio/devices/iio:device0/buffer/enable ]; then \
        echo 1 > /sys/bus/iio/devices/iio:device0/buffer/enable 2>/dev/null; \
    fi; \
    sleep 5; \
done'

[Install]
WantedBy=multi-user.target
